# build/conf/upstart.conf
# at /etc/init/asana-media.conf

description "Gunicorn application server running asana-media"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid www-data

exec /opt/asana-media/run_app.sh
